# infra/dbt/crontab
# dbt cron schedule
# Format: minute hour day month weekday command
# Default: run every hour at minute 0
0 * * * * /app/infra/dbt/run_dbt.sh run >> /app/dbt/logs/cron.log 2>&1

